# leak.service
[Unit]
Description=itterent Compose Application
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/studio
User=mrqadeer
Group=docker
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=COMPOSE_PROJECT_NAME=leak

ExecStartPre=/bin/bash -c 'timeout 30 sh -c "until docker info > /dev/null 2>&1; do sleep 1; done"'
ExecStartPre=/usr/local/bin/docker-compose down --remove-orphans
ExecStart=/usr/local/bin/docker-compose up -d --build
ExecStop=/usr/local/bin/docker-compose down --timeout 30
ExecReload=/usr/local/bin/docker-compose restart

Restart=on-failure
RestartSec=30
TimeoutStartSec=600
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target